# samba/arm32/local/etc/s6/smbd/finish
#!/bin/sh

echo "finish smbd"

[ -r /var/run/samba/smbd.pid ] && kill -9 $(cat /var/run/samba/smbd.pid)
