# samba/arm32/local/etc/s6/smbd/run
#!/bin/sh

echo "run smbd"

exec smbd -F -S -s=/conf.d/samba/smb.conf
