# samba/arm32/local/etc/s6/nmbd/finish
#!/bin/sh

echo "finish nmbd"

[ -r /var/run/samba/nmbd.pid ] && kill -9 $(cat /var/run/samba/nmbd.pid)
