# samba/arm32/local/etc/s6/nmbd/run
#!/bin/sh

echo "run nmbd"

exec nmbd -F -S -s=/conf.d/samba/smb.conf
